# chttp.socket
[Unit]
Description=A Cimple HTTP server

[Socket]
ListenStream=80
Accept=true

[Install]
WantedBy=sockets.target
